# packaging/systemd/locallens.service
[Unit]
Description=LocalLens station health check
After=local-fs.target

[Service]
Type=oneshot
User=locallens
Group=locallens
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/locallens
ExecStart=/usr/local/bin/locallens /var/lib/locallens/catalog.db health /var/lib/locallens/store
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
